# luci-app-cpufreq/root/etc/uci-defaults/21_tuning_net
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@tuning_net[-1]
	add ucitrack tuning_net
	set ucitrack.@tuning_net[-1].init=tuning_net
	commit ucitrack
EOF

exit 0